# installation/autoflower.conf
[program:autoflower]
directory=/home/pi/pio/python
command=python server.py
autostart=true
autorestart=true

